# overlays/mt3000/rootfs/etc/uci-defaults/99-mt3000-wifi-defaults
#!/bin/sh

log_tag="mt3000-wifi-defaults"
board_info="/usr/bin/mt798x-board-info"

log() {
	logger -t "$log_tag" "$*"
}

board_name() {
	cat /tmp/sysinfo/board_name 2>/dev/null
}

is_mt3000() {
	case "$(board_name)" in
		glinet,gl-mt3000|glinet,mt3000-snand)
			return 0
			;;
	esac

	return 1
}

valid_key() {
	case "$1" in
		*[!A-Za-z0-9]*|'')
			return 1
			;;
	esac

	[ ${#1} -ge 8 ] && [ ${#1} -le 63 ]
}

helper_value() {
	local cmd="$1"
	local fallback="$2"
	local value

	if [ -x "$board_info" ]; then
		value="$("$board_info" "$cmd" 2>/dev/null | tr -d '\r\n')"
	fi

	[ -n "$value" ] && {
		printf '%s\n' "$value"
		return
	}

	printf '%s\n' "$fallback"
}

derive_fallback_ssid() {
	local mac suffix

	for path in /sys/class/net/eth0/address /sys/class/net/eth1/address /sys/class/net/br-lan/address; do
		mac="$(cat "$path" 2>/dev/null | tr -d ':\r\n' | tr 'A-F' 'a-f')"
		[ ${#mac} -ge 3 ] || continue

		suffix="${mac#?????????}"
		[ -n "$suffix" ] && {
			printf 'GL-MT3000-%s\n' "$suffix"
			return
		}
	done

	printf 'GL-MT3000\n'
}

derive_fallback_key() {
	local mac

	for path in /sys/class/net/eth0/address /sys/class/net/eth1/address /sys/class/net/br-lan/address; do
		mac="$(cat "$path" 2>/dev/null | tr -d ':\r\n' | tr 'a-f' 'A-F')"
		[ ${#mac} -eq 12 ] || continue
		printf 'MT3000%s\n' "$mac"
		return
	done

	printf 'MT3000OPEN\n'
}

ensure_wireless_config() {
	[ -s /etc/config/wireless ] && return 0

	log "wireless config missing; running wifi config"
	wifi config >/tmp/mt3000-wifi-config.log 2>&1
}

radio_band() {
	local radio="$1"
	local band
	local path

	band="$(uci -q get "wireless.$radio.band")"
	[ -n "$band" ] && {
		printf '%s\n' "$band"
		return
	}

	path="$(uci -q get "wireless.$radio.path")"
	case "$radio $path" in
		*5g*|*5G*|*rax*|*MT7981_1_2*|*wifi1*|*radio1*)
			printf '5g\n'
			;;
		*)
			printf '2g\n'
			;;
	esac
}

ensure_iface() {
	local iface="$1"
	local radio="$2"
	local ssid="$3"
	local key="$4"
	local band

	[ -n "$iface" ] && [ -n "$radio" ] || return 1

	band="$(radio_band "$radio")"

	uci -q set "wireless.$iface=wifi-iface"
	uci -q set "wireless.$iface.device=$radio"
	uci -q set "wireless.$iface.mode=ap"
	uci -q set "wireless.$iface.network=lan"
	uci -q set "wireless.$iface.ssid=$ssid"
	uci -q set "wireless.$iface.disabled=0"

	if valid_key "$key"; then
		uci -q set "wireless.$iface.encryption=psk2"
		uci -q set "wireless.$iface.key=$key"
	else
		uci -q set "wireless.$iface.encryption=none"
		uci -q delete "wireless.$iface.key"
	fi

	log "configured iface=$iface radio=$radio band=$band ssid=$ssid key_present=$([ -n "$key" ] && echo yes || echo no)"
}

repair_wifi_defaults() {
	local ssid key key_source
	local radios radio iface ifaces configured

	ssid="$(helper_value wifi-ssid "$(derive_fallback_ssid)")"
	key="$(helper_value wifi-key "$(derive_fallback_key)")"
	key_source="factory_or_fallback"

	if ! valid_key "$key"; then
		key=""
		key_source="open_no_valid_key"
	fi

	radios="$(uci -q show wireless | sed -n "s/^wireless\\.\\([^.=]*\\)=wifi-device$/\\1/p")"

	for radio in $radios; do
		uci -q set "wireless.$radio.disabled=0"

		ifaces="$(uci -q show wireless | sed -n "s/^wireless\\.\\([^.=]*\\)\\.device='$radio'$/\\1/p;s/^wireless\\.\\([^.=]*\\)\\.device=$radio$/\\1/p")"
		if [ -z "$ifaces" ]; then
			iface="default_${radio}"
			ensure_iface "$iface" "$radio" "$ssid" "$key"
			configured=$((configured + 1))
			continue
		fi

		for iface in $ifaces; do
			ensure_iface "$iface" "$radio" "$ssid" "$key"
			configured=$((configured + 1))
		done
	done

	uci commit wireless
	log "committed wireless defaults ssid=$ssid key_source=$key_source iface_count=${configured:-0}"
}

is_mt3000 || exit 0

ensure_wireless_config
repair_wifi_defaults

exit 0
