# overlays/mt3000/rootfs/etc/uci-defaults/99-mt3000-wifi-defaults
#!/bin/sh

log_tag="mt3000-wifi-defaults"
board_info="/usr/bin/mt798x-board-info"

log() {
	logger -t "$log_tag" "$*"
}

board_name() {
	cat /tmp/sysinfo/board_name 2>/dev/null
}

is_mt3000() {
	case "$(board_name)" in
	glinet,gl-mt3000|glinet,mt3000-snand)
		return 0
		;;
	esac

	return 1
}

valid_key() {
	case "$1" in
	''|*[!A-Za-z0-9]*)
		return 1
		;;
	esac

	[ ${#1} -ge 8 ] && [ ${#1} -le 63 ]
}

helper_value() {
	local cmd="$1"
	local value

	if [ -x "$board_info" ]; then
		value="$("$board_info" "$cmd" 2>/dev/null | tr -d '\r\n')"
	fi

	[ -n "$value" ] && printf '%s\n' "$value"
}

fallback_ssid() {
	local mac path suffix

	for path in /sys/class/net/eth0/address /sys/class/net/eth1/address /sys/class/net/br-lan/address; do
		mac="$(cat "$path" 2>/dev/null | tr -d ':\r\n' | tr 'A-F' 'a-f')"
		[ ${#mac} -ge 3 ] || continue
		suffix="${mac#?????????}"
		printf 'GL-MT3000-%s\n' "$suffix"
		return 0
	done

	printf 'GL-MT3000\n'
}

ensure_wireless_config() {
	uci -q show wireless >/dev/null 2>&1 && return 0

	if command -v wifi >/dev/null 2>&1; then
		wifi config >/dev/null 2>&1 || true
	fi
}

radio_band() {
	local radio="$1"
	local band path phy

	band="$(uci -q get "wireless.$radio.band")"
	[ -n "$band" ] && {
		printf '%s\n' "$band"
		return 0
	}

	phy="$(uci -q get "wireless.$radio.phy")"
	for path in "/sys/class/ieee80211/$phy/name" "/sys/class/ieee80211/$phy/index"; do
		[ -r "$path" ] || continue
		cat "$path" 2>/dev/null
		return 0
	done

	printf 'unknown\n'
}

ensure_iface() {
	local iface="$1"
	local radio="$2"
	local ssid="$3"
	local key="$4"
	local band

	[ -n "$iface" ] && [ -n "$radio" ] || return 1

	band="$(radio_band "$radio")"

	uci -q set "wireless.$iface=wifi-iface"
	uci -q set "wireless.$iface.device=$radio"
	uci -q set "wireless.$iface.mode=ap"
	uci -q set "wireless.$iface.network=lan"
	uci -q set "wireless.$iface.ssid=$ssid"
	uci -q set "wireless.$iface.disabled=0"

	if valid_key "$key"; then
		uci -q set "wireless.$iface.encryption=psk2"
		uci -q set "wireless.$iface.key=$key"
	else
		uci -q set "wireless.$iface.encryption=none"
		uci -q delete "wireless.$iface.key"
	fi

	log "configured iface=$iface radio=$radio band=$band ssid=$ssid key_present=$([ -n "$key" ] && echo yes || echo no)"
}

repair_wifi_defaults() {
	local ssid key key_source radios radio iface ifaces configured

	ssid="$(helper_value ssid-2g)"
	[ -n "$ssid" ] || ssid="$(fallback_ssid)"

	key="$(helper_value wifi-key)"
	key_source="factory"

	if ! valid_key "$key"; then
		key=""
		key_source="none_open"
	fi

	radios="$(uci -q show wireless | sed -n "s/^wireless\\.\\([^.=]*\\)=wifi-device$/\\1/p")"
	configured=0

	for radio in $radios; do
		uci -q set "wireless.$radio.disabled=0"

		ifaces="$(uci -q show wireless | sed -n "s/^wireless\\.\\([^.=]*\\)\\.device='$radio'$/\\1/p;s/^wireless\\.\\([^.=]*\\)\\.device=$radio$/\\1/p")"
		if [ -z "$ifaces" ]; then
			iface="default_${radio}"
			ensure_iface "$iface" "$radio" "$ssid" "$key"
			configured=$((configured + 1))
			continue
		fi

		for iface in $ifaces; do
			ensure_iface "$iface" "$radio" "$ssid" "$key"
			configured=$((configured + 1))
		done
	done

	uci commit wireless
	log "committed wireless defaults ssid=$ssid key_source=$key_source iface_count=$configured"
}

is_mt3000 || exit 0
ensure_wireless_config
repair_wifi_defaults

exit 0
